# templates/deploy/raro-tracker-check.service
[Unit]
Description=Raro Tracker - checagem de disponibilidade (oneshot)

[Service]
Type=oneshot
WorkingDirectory=/home/rammyres/raro-tracker
EnvironmentFile=/home/rammyres/raro-tracker/.env
ExecStart=/home/rammyres/raro-tracker/venv/bin/python3 scheduler/check_availability.py
