# deploy/raro-tracker-check.service
[Unit]
Description=Raro Tracker - checagem de disponibilidade (oneshot)

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/rarotracker
EnvironmentFile=/opt/rarotracker/.env
ExecStart=/opt/rarotracker/venv/bin/python3 scheduler/check_availability.py
