# deploy/raro-tracker-check.service
[Unit]
Description=Raro Tracker - checagem de disponibilidade (oneshot)

[Service]
Type=oneshot
WorkingDirectory=/opt/rarotracker
EnvironmentFile=/opt/rarotracker/.env
ExecStart=/opt/rarotracker/check_availability.py
